# src/services/taris-web.service
[Unit]
Description=Taris Web UI (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stas
WorkingDirectory=/home/stas/.taris
ExecStart=/usr/bin/python3 /home/stas/.taris/bot_web.py
Restart=on-failure
RestartSec=10
MemoryMax=400M
MemoryHigh=350M
KillMode=mixed
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/stas
EnvironmentFile=/home/stas/.taris/bot.env
Environment=TARIS_BIN=/usr/bin/picoclaw

[Install]
WantedBy=default.target
