# deploy/system-configs/vps/code-server.service
[Unit]
Description=code-server (browser VS Code)
After=network.target

[Service]
Type=simple
User=stas
Environment=HOME=/home/stas
# GITHUB_TOKEN is set from /home/stas/.taris-env or injected at deploy time — never commit the value
# Environment=GITHUB_TOKEN=<token>
EnvironmentFile=-/home/stas/.code-server-env
ExecStart=/usr/local/bin/code-server --config /home/stas/.config/code-server/config.yaml /home/stas
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
