# deploy/system-configs/vps/ttyd.service
[Unit]
Description=ttyd - Web Terminal
After=network.target

[Service]
User=stas
Group=stas
ExecStart=/usr/bin/ttyd --port 7681 --interface 127.0.0.1 /bin/bash -l
Restart=always
RestartSec=5
Environment=TERM=xterm-256color
Environment=HOME=/home/stas
WorkingDirectory=/home/stas

[Install]
WantedBy=multi-user.target
