# deploy/system-configs/taristation2/systemd/user/taris-web.service
[Unit]
Description=Taris Web UI (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/stas/.taris
ExecStart=/usr/bin/python3 /home/stas/.taris/bot_web.py
Restart=on-failure
RestartSec=10
MemoryMax=400M
MemoryHigh=350M
KillMode=mixed
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/stas/.taris/bot.env

[Install]
WantedBy=default.target
